# rustdesk-server/rootfs/etc/s6-overlay/s6-rc.d/hbbr/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant App: RustDesk-Server
# Runs hbbr
# ==============================================================================


cd /config
PARAMS=
[ "${ENCRYPTED_ONLY}" = "1" ] && PARAMS="-k _"


bashio::log.info $(/opt/rustdesk-server/hbbr --version)
/opt/rustdesk-server/hbbr $PARAMS
